# deploy/technews-backend.service
[Unit]
Description=Tech News Chatbot — FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data

WorkingDirectory=/var/www/technews/backend

# Loads all variables from .env into the process environment
EnvironmentFile=/var/www/technews/backend/.env

# Use the venv's uvicorn; bind to loopback only (nginx is the public face)
ExecStart=/var/www/technews/backend/.venv/bin/uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1

Restart=always
RestartSec=5

# Keep logs tidy
StandardOutput=journal
StandardError=journal
SyslogIdentifier=technews-backend

[Install]
WantedBy=multi-user.target
